# trading/bot/dashboard/msd/msd-dashboard.service
[Unit]
Description=MSD Web Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/clawd-biz/trading/bot/dashboard/msd
ExecStart=/home/ubuntu/clawd-biz/trading/.venv/bin/python app_msd.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
